# deploy/systemd/neo-grace.timer
[Unit]
Description=Daily subscription grace/expiry reminders

[Timer]
OnCalendar=daily
Persistent=true
Unit=neo-grace.service

[Install]
WantedBy=timers.target
