# deploy/systemd/neo-rollup.timer
[Unit]
Description=Recompute sales rollups hourly

[Timer]
OnCalendar=hourly
Persistent=true
Unit=neo-rollup.service

[Install]
WantedBy=timers.target
