# deploy/systemd/neo-anonymize.timer
[Unit]
Description=Run tenant PII anonymization daily

[Timer]
OnCalendar=daily
Persistent=true
Unit=neo-anonymize.service

[Install]
WantedBy=timers.target
